# templates/default/systemd.service.erb
[Unit]
Description="Kibana <%= node['kibana']['version'] %>"

[Service]
User=<%= @user %>
Restart=always
ExecStart=<%= File.join(@base_dir, @bin_path, 'kibana') %>

[Install]
WantedBy=multi-user.target
