# templates/default/upstart.conf.erb
# !!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!
# !! Managed by Chef, do not edit !!
# !!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!
description "Kibana Server"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

<% if @recent_upstart %>
setuid <%= node['kibana']['user'] %>
setgid <%= node['kibana']['group'] %>
<% end %>
respawn limit 15 5
umask 0077

pre-start script
  test -d <%= node['kibana']['base_dir'] %> || { stop; exit 0; }
end script

chdir <%= node['kibana']['base_dir'] %>
<%
  command = case @version
            when 3
              raise "kibana 3 cannot be started as a service"
            when 4
              "#{@bin_path}/kibana #{@options}"
            end
  if node['kibana']['port'] < 1024
    command = "/usr/bin/authbind --deep " + command
  end
%>

<% if @recent_upstart %>
exec <%= command %>
<% else %>
exec su -s /bin/sh -c 'exec "$0" "$@"' <%= node['kibana']['user'] %> -- <%= command %>
<% end %>
