# templates/default/upstart.conf.erb
#
# Generated by Chef
#
description     "Kibana Server"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

setuid <%= node['kibana']['user'] %>
setgid <%= node['kibana']['group'] %>
respawn limit 15 5
umask 0077

pre-start script
  test -d <%= node['kibana']['base_dir'] %> || { stop; exit 0; }
end script

chdir <%= node['kibana']['base_dir'] %>
exec <%= node['kibana']['port'] < 1024 ? "/usr/bin/authbind --deep " : "" %>ruby1.9.1 kibana.rb
